# overlay/etc/supervisor/conf.d/cozy-init.conf
[program:cozy-init]
autorestart=false
autostart=true
command=/usr/bin/nohup /etc/init.d/cozy-init
redirect_stderr=true
user=root
priority=1
