# patches/etc/supervisor/conf.d/couchdb.conf
[program:couchdb]
autorestart=true
autostart=true
command=couchdb
environment=HOME=/var/lib/couchdb
directory=/var/lib/couchdb
redirect_stderr=true
user=couchdb
priority=1
